# dist/tpm2-totp.service.in
[Unit]
Description=Display a TOTP during boot
Requires=systemd-vconsole-setup.service dev-tpm0.device
After=systemd-vconsole-setup.service dev-tpm0.device
Before=systemd-ask-password-console.service
Conflicts=multi-user.target
DefaultDependencies=no

[Service]
Environment="TPM2TOTP_TCTI=device:/dev/tpm0"
Type=exec
ExecStart=@HELPERSDIR@/show-tpm2-totp
StandardOutput=tty

[Install]
WantedBy=sysinit.target
